# neutron-agents/supervisord.conf
[program:neutron-linuxbridge-agent]
command=/opt/kolla/neutron/config-linuxbridge-agent.sh
priority=30
startsec=10
autostart=true
stderr_events_enabled=true
stdout_events_enabled=true
